# files/etc/systemd/system/firstboot-brew-install.service
[Unit]
Description=Install Homebrew for user

[Service]
Type=oneshot
ExecStart=/usr/local/bin/install-homebrew.sh %i
