# files/usr/lib/systemd/user/firstboot-setup.service
[Unit]
Description=Run First Boot Setup Wizard
After=xdg-desktop-autostart.target
Wants=xdg-desktop-autostart.target
ConditionPathExists=%h/.config/firstboot/run

[Service]
Type=exec
ExecStart=%h/.config/firstboot/firstboot-setup.sh

[Install]
WantedBy=xdg-desktop-autostart.target

