# files/usr/lib/systemd/user/firstboot-setup.service
[Unit]
Description=Run First Boot Setup Wizard
After=graphical-session.target
Wants=graphical-session.target
ConditionPathExists=%h/.config/firstboot/run

[Service]
Type=exec
ExecStart=%h/.config/firstboot/firstboot-setup.sh

[Install]
WantedBy=graphical-session.target
