# files/usr/lib/systemd/user/firstboot-setup.service
[Unit]
Description=First Boot Wizard
After=graphical-session.target cosmic-session.target
Wants=graphical-session.target cosmic-session.target
ConditionPathExists=%h/.config/cosmic-initial-setup.done
ConditionPathExists=%h/.config/firstboot/run

[Service]
Type=simple
ExecStart=%h/.config/firstboot/firstboot-setup.sh

[Install]
WantedBy=default.target

