# files/usr/lib/systemd/user/firstboot-setup.service
[Unit]
Description=First Boot Wizard
ConditionPathExists=%h/.config/firstboot/run

[Service]
Type=simple
ExecStart=%h/.config/firstboot/firstboot-setup.sh

[Install]
WantedBy=default.target
