# files/usr/lib/systemd/system/firstboot-cleanup.service
[Unit]
Description=Remove firstboot marker

[Service]
Type=oneshot
ExecStart=/usr/bin/rm -f /etc/myos/firstboot
